# scripts/arch/usr/local/bin/optimization/fstrim.service
[Unit]
Description=fstrim_sdcard_optimization

[Service]
Type=oneshot
ExecStart=/usr/local/bin/optimization/fstrim.sh
