# scripts/arch/usr/local/bin/duckdns/duckdns.service
[Unit]
Description=dynamic_dns_script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns/duckdns.sh
